# luci/applications/luci-DNS-Filter/root/etc/crontabs/root
00 0 * * * screen -S DNS -dm bash dns update 2>/dev/null 2>/dev/null #DNS-Filter 
