# packages/ddns-openrouter/files/etc/init.d/ddns
#!/bin/sh /etc/rc.common

START=95

start () {
   touch /tmp/limeng_ddns_start.txt
   sh /usr/lib/ddns/dynamic_dns_updater.sh start 
}

stop () {

  echo "ddns stop"
}

